# dotfiles/config/autostart/batcheck.desktop
[Desktop Entry]
Name=battery checker
Comment=daemon which periodically check the battery status and warns the user
Type=Application
Exec=batcheck
