# dotfiles/config/autostart/browser.desktop
[Desktop Entry]
Name=browser
Comment=Starts the default browser, defined by XDG.
Type=Application
Exec=browser
