# dotfiles/config/autostart/irc.desktop
[Desktop Entry]
Name=IRC
Comment=SSH session to your IRC server
Type=Application
Exec=i3-sensible-terminal -e ssh irc
