# dotfiles/config/autostart/irc.desktop
[Desktop Entry]
Name=IRC
Comment=Starts a mosh session to your IRC server.
Type=Application
Exec=urxvtc -e mosh irc
