# dotfiles/config/autostart/irc.desktop
[Desktop Entry]
Name=IRC
Comment=mosh session to your IRC server
Type=Application
Exec=urxvtc -name irc -e mosh irc
