# dotfiles/config/autostart/nobell.desktop
[Desktop Entry]
Name=Disable X bell
Comment=Disables the X bell
Type=Application
Exec=xset b off
