# dotfiles/config/autostart/sc.desktop
[Desktop Entry]
Name=screen config
Comment=Automatically changes mode and position of the secondary screen output.
Type=Application
Exec=sc auto
