# dotfiles/config/autostart/terminal.desktop
[Desktop Entry]
Name=Terminal
Comment=An available terminal emulator
Type=Application
Exec=i3-sensible-terminal
