# dotfiles/config/autostart/unclutter.desktop
[Desktop Entry]
Name=unclutter
Comment=Hides incative mouse pointer
Type=Application
Exec=unclutter -noevents
