# dotfiles/config/autostart/urxvtc.desktop
[Desktop Entry]
Name=rxvt-unicode-client
Comment=rxvt controll program
Type=Application
Exec=urxvtc
