# dotfiles/config/autostart/xkbcfg.desktop
[Desktop Entry]
Name=xkbcfg
Comment=Modifies keymaps and pointer button mappings in this X session.
Type=Application
Exec=xkbcfg
